# provision/roles/dokku-leak-stop/files/dokku-leak-stop.timer
[Unit]
Description=Stop dokku leak processes periodically

[Timer]
OnBootSec=30min
OnUnitActiveSec=5h 30min

[Install]
WantedBy=timers.target
