# provision/roles/sync-hwclock/files/sync-hwclock.service
[Unit]
Description=Sync hwclock

[Service]
Type=oneshot
ExecStartPre=/sbin/hwclock --show
ExecStart=/sbin/hwclock --systohc
