# provision/roles/sync-hwclock/files/sync-hwclock.timer
[Unit]
Description=Sync hwclock periodically

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
